# vpnclient/lede/bin-lede-4099-g4c3953ba29-airrouter-futel-work-in-progress/files/etc/uci-defaults/futel.zzz.defaults
#!/bin/sh

uci show > /etc/uci-firstboot
